# src/main/resources/cron.d/temperature
* * * * *  root bash -c "/usr/bin/python /opt/weather/src/main/resources/temperature.py > /var/log/temperature.log 2>/var/log/cron.log"
